# ansible_vm_setup/roles/app/templates/medusa-backend.service.j2
[Unit]
Description=Medusa v2 Backend API
After=network.target

[Service]
Type=simple
User=medusa
WorkingDirectory=/opt/medusa-backend
EnvironmentFile=/opt/medusa-backend/.env
ExecStart=/usr/bin/node .medusa/server/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medusa-backend
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
