# ansible_vm_setup/roles/web/templates/medusa-storefront.service.j2
[Unit]
Description=Medusa Next.js Storefront
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/opt/medusa-storefront
EnvironmentFile=/opt/medusa-storefront/.env.local
ExecStart=/usr/bin/node node_modules/.bin/next start -p {{ storefront_port }}
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=medusa-storefront

[Install]
WantedBy=multi-user.target
